# files/dap2965/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Never uncomment this unless YOU KNOW WHAT YOU'RE DOING
#for file in /sys/kernel/debug/ieee80211/*/ath10k/dfs_block_radar_events; do
#    echo Y > $file
#done

/usr/sbin/ebtables -A FORWARD -i wlan+ -o wlan+ -j DROP
/usr/sbin/ebtables -A FORWARD -i wlan+ -s 00:1e:0b:d9:bb:6a -j DROP
/usr/sbin/ebtables -A FORWARD -i eth+ -o wlan+ -s ! 00:1e:0b:d9:bb:6a -j DROP

exit 0
